# 80-xmm7360.rules
# /usr/lib/udev/rules.d/80-xmm7360.rules

# Prevent PCIe runtime power management from auto-suspending the modem.
ACTION=="add", SUBSYSTEM=="pci", \
    ATTR{vendor}=="0x8086", ATTR{device}=="0x7360", \
    ATTR{power/control}="on"

# Tag ttyXMM AT ports for the Fibocom MM plugin (signal, registration)
SUBSYSTEM=="tty", KERNEL=="ttyXMM*", ENV{ID_MM_FIBOCOM_TAGGED}="1"

# Trigger xmm7360-init.service when the first AT port appears.
# Use restart (not SYSTEMD_WANTS) so it re-runs on every fresh tty-add —
# a oneshot+RemainAfterExit service would otherwise refuse to re-run
# after a manual modprobe, leaving MM without the RPC wake.
SUBSYSTEM=="tty", KERNEL=="ttyXMM1", ACTION=="add", \
    TAG+="systemd", \
    RUN+="/usr/bin/systemctl restart --no-block xmm7360-init.service"

# ── Kernel-driven recovery via uevent ────────────────────────────────────
# The module emits CHANGE uevents with XMM7360_STATE=<state>.
# Two states matter:
#   failed  → kernel recovery exhausted; need module reload
#   resumed → woke from suspend; re-run RPC init

ACTION=="change", SUBSYSTEM=="pci", \
    ATTR{vendor}=="0x8086", ATTR{device}=="0x7360", \
    ENV{XMM7360_STATE}=="failed", \
    RUN+="/usr/bin/systemctl --no-block start xmm7360-recovery.service"

ACTION=="change", SUBSYSTEM=="pci", \
    ATTR{vendor}=="0x8086", ATTR{device}=="0x7360", \
    ENV{XMM7360_STATE}=="resumed", \
    RUN+="/usr/bin/systemctl restart xmm7360-init.service"

# After kernel-level recovery rebuilds the devices, re-run RPC init
# and force ModemManager to re-scan (it lost the old modem object).
ACTION=="change", SUBSYSTEM=="pci", \
    ATTR{vendor}=="0x8086", ATTR{device}=="0x7360", \
    ENV{XMM7360_STATE}=="recovered", \
    RUN+="/usr/bin/systemctl restart xmm7360-init.service", \
    RUN+="/usr/bin/mmcli -S"
