# xmm7360-init.service
[Unit]
Description=XMM7360 Modem RPC Initialisation
After=dev-ttyXMM1.device
Before=ModemManager.service

[Service]
Type=oneshot
RemainAfterExit=yes

EnvironmentFile=-/etc/xmm7360.conf

ExecStart=/usr/bin/open_xdatachannel --init-only

ExecStartPost=/usr/bin/touch /run/xmm7360-init.done

ExecStop=/usr/bin/rm -f /run/xmm7360-init.done

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
