# xmm7360-resume.service
[Unit]
Description=XMM7360 Modem Post-Resume Recovery
# Run after system is fully back up and MM is ready
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
After=ModemManager.service
Requires=ModemManager.service

[Service]
Type=oneshot
# --no-block so this service returns immediately and recovery runs in background
ExecStart=/usr/bin/systemctl start --no-block xmm7360-recovery.service

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
