# xmm7360-cleanup.service
[Unit]
Description=XMM7360 post-disconnect modem cycle
# Triggered by udev when the kernel emits XMM7360_STATE=disconnected
# after a PPP session tears down. A quick AT+CFUN cycle (via mmcli over
# D-Bus) clears the firmware's stale data-session state so the next
# connect succeeds. Without this, manual disconnect -> reconnect wedges
# the modem until a module reload.

[Service]
Type=oneshot
RemainAfterExit=no
TimeoutStartSec=30

ExecStart=/bin/sh -c '\
    logger -t xmm7360-cleanup "post-disconnect modem cycle"; \
    m=$(mmcli -L 2>/dev/null | grep -oE "/org/freedesktop/ModemManager1/Modem/[0-9]+" | head -1); \
    if [ -z "$m" ]; then \
        logger -t xmm7360-cleanup "no modem present, skipping"; \
        exit 0; \
    fi; \
    timeout 10 mmcli -m "$m" --disable 2>/dev/null || true; \
    sleep 1; \
    timeout 15 mmcli -m "$m" --enable  2>/dev/null || true; \
    logger -t xmm7360-cleanup "cycle complete"'

StandardOutput=journal
StandardError=journal
