# systemd/xmm7360-signal.service
[Unit]
Description=XMM7360 signal-quality polling setup
# Runs AFTER ModemManager is up (unlike init.service which wakes the radio
# before MM). Enables periodic signal reporting once MM has the modem.
After=ModemManager.service xmm7360-init.service
Wants=ModemManager.service
ConditionPathExists=/run/xmm7360-init.done

[Service]
Type=oneshot
RemainAfterExit=yes

# Wait for MM to actually have the modem, then enable signal polling.
# Resolve the modem index dynamically: on a multi-modem system (or after a
# recovery re-enumeration) the XMM7360 is not necessarily index 0, so match
# it by its ttyXMM ports / Fibocom identity rather than hardcoding -m 0.
ExecStart=/bin/sh -c '\
    for i in $(seq 1 30); do \
        idx=""; \
        for p in $(mmcli -L 2>/dev/null | grep -oE "/Modem/[0-9]+"); do \
            n=${p##*/}; \
            if mmcli -m "$n" 2>/dev/null | grep -qiE "ttyXMM|XMM7360|Fibocom|L850"; then \
                idx="$n"; break; \
            fi; \
        done; \
        if [ -n "$idx" ] && mmcli -m "$idx" --signal-setup=5 2>/dev/null; then \
            exit 0; \
        fi; \
        sleep 1; \
    done; \
    exit 0'

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
